# sv/tmp/run
#!/bin/sh
exec chpst -uroot /usr/bin/main 2>&1
